# systemd/discord-rss-bot.service
[Unit]
Description=Discord RSS Bot
Documentation=https://github.com/isyuricunha/discord-news-rss-bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
DynamicUser=yes
StateDirectory=discord-rss-bot
StateDirectoryMode=0750
Environment=RSS_BOT_DATA=/var/lib/discord-rss-bot
Environment=RSS_BOT_CONFIG=/etc/discord-rss-bot/config.env
EnvironmentFile=-/etc/discord-rss-bot/config.env
ExecStart=/usr/local/bin/discord-rss-bot run
Restart=on-failure
RestartSec=10s
StartLimitIntervalSec=5m
StartLimitBurst=5
KillSignal=SIGTERM
TimeoutStopSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-rss-bot

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/discord-rss-bot
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
