# systemd/discord-rss-bot.service
[Unit]
Description=Discord RSS Bot Service
Documentation=https://github.com/YOUR_USERNAME/discord-rss-bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=discord-rss-bot
Group=discord-rss-bot
WorkingDirectory=/opt/discord-rss-bot
ExecStart=/usr/bin/python3 /opt/discord-rss-bot/bot_service.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-rss-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/discord-rss-bot /var/log/discord-rss-bot
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Environment
Environment=RSS_BOT_DATA=/var/lib/discord-rss-bot
Environment=RSS_BOT_LOGS=/var/log/discord-rss-bot
Environment=RSS_BOT_CONFIG=/etc/discord-rss-bot/config.env
Environment=RSS_BOT_PID=/var/run/discord-rss-bot.pid

[Install]
WantedBy=multi-user.target
